# scripts/backup-db.timer
[Unit]
Description=Run the ImpositionPDF SQLite backup daily

[Timer]
# Every day at 03:30 local time, catching up if the machine was off.
OnCalendar=*-*-* 03:30:00
Persistent=true

[Install]
WantedBy=timers.target
